# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:app1]
directory=/app1
command=java -jar /usr/src/HTTPTrafficFakeLog.jar -v
autorestart=true

[program:app2]
directory=/app2
command=java -jar /usr/src/HTTPLogMonitoringTool.jar
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true